# debian/postrm
#!/bin/bash

set -x

systemctl disable m51-lighthouse
systemctl daemon-reload
systemctl stop m51-lighthouse
